# systemd/torbox-rclone-mount.service
[Unit]
Description=Cloudarr WebDAV rclone mount
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=cloudarr
Group=cloudarr
ExecStart=/usr/bin/rclone mount torboxdav: /mnt/torbox --config /etc/rclone/rclone.conf --allow-other --dir-cache-time 30s --vfs-cache-mode off --poll-interval 15s --rc --rc-no-auth
ExecStop=/bin/fusermount -u /mnt/torbox
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
